# grafana-server/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network-online.target
[Service]
Restart=on-failure
ExecStart=/usr/local/bin/node_exporter
[Install]
WantedBy=multi-user.target
